# s6/pre-migrate/run
#!/command/execlineb -P
# vim: set syntax=sh :

python3 -m alembic upgrade head
