# s6/pre-migrate/up
/bin/sh -c "cd /app && python3 -m alembic upgrade head"